# Config/files/etc/crontabs/root
20 5 8 * * sleep 5 && touch /etc/banner && reboot
*/5 * * * * /etc/wan_check.sh
