# Config/files/etc/crontabs/root
#每月定时重启路由器
11 1 11 * * sleep 5 && touch /etc/banner && reboot
#每5分钟检测网络，重启WAN
*/5 * * * * /etc/install_ipk.sh
#############################
