# Config/files/etc/crontabs/root
11 1 11 * * sleep 5 && touch /etc/banner && reboot
*/5 * * * * /etc/wan_check.sh
